# .config/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Comments=Simple X hotkey daemon
Exec=/usr/bin/sxhkd
Termianl=false
Type=Application